# deploy/csat-update.service
[Unit]
Description=CSAT auto-update (latest release within the tracked major)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/csat-update
